# etc/init/telega-web.conf
start on syno.network.ready or runlevel 2
stop on runlevel [06]

respawn
respawn limit 5 10

mkdir -p /var/log/telega/
exec /usr/bin/telega-web &>/var/log/telega/web.log
