# etc/init/telega-worker.conf
start on syno.network.ready or runlevel 2
stop on runlevel [06]

respawn
respawn limit 5 10

exec /usr/bin/telega-worker
